# mcp-server/systemd/career-mcp.service
[Unit]
Description=Career Context Server for Job Match AI
After=network.target

[Service]
Type=simple
User=volshin
WorkingDirectory=/home/volshin/job-match-ai/mcp-server
Environment="PATH=/home/volshin/job-match-ai/mcp-server/venv/bin"
ExecStart=/home/volshin/job-match-ai/mcp-server/venv/bin/python career_mcp.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
